# oscam/files/oscam.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2017 OpenWrt.org


START=99


start() {
	/usr/bin/oscam -b -u
	echo "runing!"
}

stop() {
	killall oscam
	killall -9 oscam 
	echo "stop!"
}
